# deploy/com.bud.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bud.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOU/bud2/bin/bud</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/YOU/bud2</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/YOU</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/YOU/Library/Logs/bud.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/YOU/Library/Logs/bud.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
